# etc/supervisor/conf.d/ulordrig.conf
[program:ulordrig]
command=/ulordrig/build/ulordrig -l /var/log/ulordrig.log
autostart = true
autorestart = true
priority = 20